# AP_intercom_agent/rootfs/etc/services.d/intercom_agent/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

#
USERNAME="$(bashio::config 'username')"
PASSWORD="$(bashio::config 'password')"
PROXY="$(bashio::config 'proxy')"
REGISTER_URI="$(bashio::config 'register_uri')"
MESSAGE_URI="$(bashio::config 'message_uri')"
ENTRANCE_URI="$(bashio::config 'entrance_uri')"
ALLOWED_CALLER="$(bashio::config 'allowed_caller')"

ANSWER_CALLS="$(bashio::config 'answer_calls')"
SEND_MESSAGES="$(bashio::config 'send_messages')"

EXPIRY="600s"
RETRY="30s"

#
bashio::log.info "Starting intercom_agent (user=${USERNAME}, proxy=${PROXY}, allowed=${ALLOWED_CALLER}, answer_calls=${ANSWER_CALLS}, send_messages=${SEND_MESSAGES}, expiry={$EXPIRY}, retry={$RETRY})"

exec /usr/local/bin/intercom_agent \
  -expiry "${EXPIRY}" \
  -retry "${RETRY}" \
  -username "${USERNAME}" \
  -password "${PASSWORD}" \
  -proxy "${PROXY}" \
  -allowed-caller "${ALLOWED_CALLER}" \
  -register-uri "${REGISTER_URI}" \
  -message-uri "${MESSAGE_URI}" \
  -entrance-uri "${ENTRANCE_URI}" \
  -answer-calls="${ANSWER_CALLS}" \
  -send-messages="${SEND_MESSAGES}"
